# deploy/meteo-france.service
[Unit]
Description=Météo-France forecast map generator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/opt/meteo-map-server
ExecStart=/opt/meteo-map-server/run_all.sh
StandardOutput=journal
StandardError=journal
